# _systems/backup/etc/systemd/system/certbot-renew.service
[Unit]
Description=Renew TLS certs with Certbot

[Service]
Type=oneshot
ExecStart=/usr/bin/certbot renew
